# debian/crowdsec-firewall-bouncer-iptables.prerm
systemctl stop crowdsec-firewall-bouncer-oht || echo "cannot stop service"
systemctl disable crowdsec-firewall-bouncer-oht || echo "cannot disable service"
